# init/happn-cluster-conductor-server.conf
start on runlevel [2345]
stop on runlevel [016]

respawn

chdir /home/happn/happn-cluster-conductor
setuid happn
setgid happn

exec bin/server # logs to /var/log/upstart/..
